# coffeehouse.service
[Unit]
Description=Coffeehouse POS System
After=network.target

[Service]
Type=simple
User=nodejs
WorkingDirectory=/var/www/coffeehouse
Environment=NODE_ENV=production
Environment=PORT=7080
ExecStart=/var/www/coffeehouse/start.sh
ExecStop=/bin/kill -SIGTERM $MAINPID
ExecReload=/bin/kill -SIGHUP $MAINPID
Restart=on-failure
RestartSec=10s
LimitNOFILE=65536
LimitNPROC=65536
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
